# backend/waitlist-api.service
[Unit]
Description=Prosody Labs Waitlist API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/prosodylabs/backend
Environment=NODE_ENV=production
EnvironmentFile=/opt/prosodylabs/backend/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
